# root/etc/services.d/php-fpm/run
#!/usr/bin/execlineb -P
/usr/sbin/php-fpm -Fc /etc/php