# root/etc/s6-overlay/s6-rc.d/crond/run
#!/command/execlineb -P
crond -fL /dev/null