# root/etc/services.d/crond/run
#!/usr/bin/execlineb -P
/usr/sbin/crond -fL /dev/null